# api/supervisord.ini
[supervisord]
nodaemon=true
environment=PYTHONPATH="/api"

[program:alembic]
command=alembic upgrade head
directory=/api
autostart=true
autorestart=false
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:api]
command=uvicorn main:app --host 0.0.0.0 --port 8000 --reload
directory=/api
autostart=true
autorestart=true
depends_on=alembic
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
environment=PYTHONUNBUFFERED=1,PYTHONPATH="/api"

[program:kafka_producer]
command=python kafka_producer/main.py
directory=/api
autostart=true
autorestart=true
depends_on=api
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
environment=PYTHONUNBUFFERED=1,PYTHONPATH="/api"
